# build/rpi/content-tpl/etc/init/rpi-bluetooth.conf
start on filesystem and stopped udevtrigger

respawn
respawn limit 3 30

script
    command -v rpi_revision >/dev/null && rpi_revision | grep -qE "BCM2837|BCM2838|Model Zero W" || { stop; exit 0; }

    sed 's/[\t ]//g' /boot/config.txt | grep -q ^'dtoverlay=pi3-disable-bt' || [ ! -e /proc/device-tree/soc/gpio@7e200000/bt_pins ] && { stop; exit 0; }

    sed 's/[\t ]//g' /boot/config.txt | grep -q '^dtoverlay=pi3-miniuart-bt' && { DEVICE=ttyS0; SPEED=460800; } || { DEVICE=ttyAMA0; SPEED=921600; }

    for i in $(seq 1 5); do
        [ -e /dev/$DEVICE ] && break || echo "/dev/$DEVICE still not there, waiting .."
        sleep 1
    done

    [ ! -e /dev/$DEVICE ] && { stop; exit 0; }

    echo "using $DEVICE@$SPEED"
    exec hciattach -n /dev/$DEVICE bcm43xx $SPEED noflow -
end script

post-start script
    for i in $(seq 1 5); do
        sleep 1
        if [ ! -z "$(ls /sys/class/bluetooth 2>/dev/null)" ]; then
            [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=1/' /etc/default/bluetooth
            for d in $(ls /sys/class/bluetooth 2>/dev/null); do
                hciconfig $d up || :
            done
            break
        fi
    done
    if [ -z "$(ls /sys/class/bluetooth 2>/dev/null)" ]; then
        [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=0/' /etc/default/bluetooth
    fi
end script

post-stop script
    [ $(xbian-arch) = RPI ] || exit 0

    if [ -z "$(ls /sys/class/bluetooth 2>/dev/null)" ]; then
        [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=0/' /etc/default/bluetooth
    else
        [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=1/' /etc/default/bluetooth
    fi
end script
