# build/rpi-3.17.y/content-tpl/etc/init/rpi-bluetooth.conf
start on filesystem and stopped udevtrigger

script
    [ $(xbian-arch) = RPI ] || exit 0

    command -v rpi_revision > /dev/null && rpi_revision | grep -q BCM2837 || exit 0

    sed 's/[\t ]//g' /boot/config.txt | grep -v '^#' | grep -q 'dtoverlay=pi3-disable-bt' && exit 0
    [ -e /proc/device-tree/soc/gpio@7e200000/bt_pins ] || exit 0

    #echo "*** rpi-bluetooth start ***"
    [ -e /dev/ttyAMA0 ] || { echo "/dev/ttyAMA0 not there, waiting 5s"; sleep 5; }

    sed 's/[\t ]//g' /boot/config.txt | grep -v '^#' | grep -q 'dtoverlay=pi3-miniuart-bt' && DEVICE=ttyS0 || DEVICE=ttyAMA0
    echo "using $DEVICE"
    /usr/bin/hciattach /dev/$DEVICE bcm43xx 921600 noflow - || :
end script

post-start script
    [ $(xbian-arch) = RPI ] || exit 0

    if [ -z "$(ls /sys/class/bluetooth 2>/dev/null)" ]; then
        [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=0/' /etc/default/bluetooth
    else
        [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=1/' /etc/default/bluetooth
        for d in $(ls /sys/class/bluetooth 2>/dev/null); do
            hciconfig $d up || :
        done
    fi

end script

post-stop script
    [ $(xbian-arch) = RPI ] || exit 0

    if [ -z "$(ls /sys/class/bluetooth 2>/dev/null)" ]; then
        [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=0/' /etc/default/bluetooth
    else
        [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=1/' /etc/default/bluetooth
    fi
end script

