# build/rpi-3.17.y/content-tpl/etc/init/rpi3-bluetooth.conf
start on stopped udevtrigger and stopped module-init-tools

#
# TODO: (better) detection if Raspberry Pi 3 is used
#

script
    [ $(xbian-arch) = RPI ] || exit 0

    echo "*** rpi3-bluetooth start ***"
    [ -e /dev/ttyAMA0 ] || { echo "/dev/ttyAMA0 not there, waiting 5s"; sleep 5; }
    [ -e /proc/device-tree/soc/gpio@7e200000/bt_pins ] || echo "/proc/device-tree/soc/gpio@7e200000/bt_pins not there, isn't this RPi3?"

    /usr/bin/hciattach /dev/ttyAMA0 bcm43xx 921600 noflow - || :

end script

post-start script
    [ $(xbian-arch) = RPI ] || exit 0

    if [ -z "$(ls /sys/class/bluetooth 2>/dev/null)" ]; then
        [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=0/' /etc/default/bluetooth
        for d in $(ls /sys/class/bluetooth 2>/dev/null); do
            hciconfig $d up || :
        done
    else
        [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=1/' /etc/default/bluetooth
    fi

end script

post-stop script
    [ $(xbian-arch) = RPI ] || exit 0

    if [ -z "$(ls /sys/class/bluetooth 2>/dev/null)" ]; then
        [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=0/' /etc/default/bluetooth
    else
        [ ! -e /etc/default/bluetooth ] || sed -i 's/BLUETOOTH_ENABLED=./BLUETOOTH_ENABLED=1/' /etc/default/bluetooth
    fi
end script

