# content/etc/init/urandom.conf
# upstart imitation of urandom for sysv compatibility, part of frandom

start on started frandom
stop on stopped frandom

