# build/content-tpl/lib/udev/rules.d/01-frandom.rules
KERNEL=="random", RUN+="/usr/local/sbin/xbian-frandom"
KERNEL=="[feu]random", RUN+="/usr/local/sbin/xbian-frandom"
