# content/lib/udev/rules.d/01-hwrng.rules
KERNEL=="random", RUN+="/usr/local/sbin/xbian-hwrng"
KERNEL=="urandom", RUN+="/usr/local/sbin/xbian-hwrng"
KERNEL=="hw_random", RUN+="/usr/local/sbin/xbian-hwrng"
