# build/imx6/content-tpl/DEBIAN/postinst
#!/bin/bash

KVER="__KVER__"
KVER=$(echo $KVER | tr '~' '-')
export KVER

[ -e /lib/libarmmem.so ] && { grep -q 'libarmmem.so' /etc/ld.so.preload 2>/dev/null|| echo '/lib/libarmmem.so' >> /etc/ld.so.preload; }
if [ "$(stat -c %d:%i /)" != "$(stat -c %d:%i /proc/1/root/. 2>/dev/null)" ]; then
        [ -e /etc/ld.so.preload ] && mv /etc/ld.so.preload /etc/ld.so.preload.chroot || :
fi

if [ $1 == "configure" ]; then

        [ -h /lib/modules/$KVER.xbian/build ] && mv /lib/modules/$KVER.xbian/build /lib/modules/$KVER
        [ -d /lib/modules/$KVER.xbian ] && rm -fr /lib/modules/$KVER.xbian || true

        insserv -fr rng-tools &>/dev/null
        udevadm control --reload-rules || :

        if [ "$(/usr/local/sbin/xbian-arch)" = iMX6 ]; then
            [ ! -e /boot/mks ] || { rm -f /usr/local/sbin/mks; ln -s /boot/mks /usr/local/sbin/; :; }

            if [ -b /dev/mmcblk0 ]; then
                echo "Updating u-boot"
                dd if=/etc/uboot-env/SPL of=/dev/mmcblk0 bs=1K seek=1
                dd if=/etc/uboot-env/u-boot.img of=/dev/mmcblk0 bs=1K seek=42 conv=fsync
                if [ ! -e /boot/noenv ]; then
                    uboot-env del -i || uboot-env del -I
                    uboot-env set < /etc/uboot-env/default.txt
                    uboot-env set script boot.scr
                    uboot-env set bootdelay 2
                fi
            fi

            [ -e /dev/root ] && rr=$(readlink -e /dev/root); [ -z "$rr" ] && rr=$(findmnt -o source / -n)
            [ -e $rr ] && fstype="$(blkid -o value -s TYPE $rr)" || fstype="$(findmnt -n -r -o FSTYPE /)"
            case $fstype in
                btrfs)
                    ;;
                zfs)
                    fstype=other
                    rr="ZFS=$rr"
                    ;;
                *)
                    fstype=other
                    ;;
            esac
            sed -i "s%setenv fstype.*%setenv fstype $fstype%" /boot/boot.scr.txt

            ccl="setenv customcmdline"

            [ -n "$rr" ] && sed -i "s%root=/dev/mmcblk0p2%root=$rr%" /boot/boot.scr.txt
            if [ -e /boot/boot.scr.txt.user ]; then
                ### keep video= from previous conf
                res=$(grep video=mxcfb0:dev= /boot/boot.scr.txt.user); res=${res##*video=mxcfb0\:dev=}; res=${res%% *};
                sed -i "s%hdmi,1920x1080.*@60,if=RGB24%$res%" /boot/boot.scr.txt

                ### keep customcmdline from previous conf
                ccl=$(grep 'setenv customcmdline' /boot/boot.scr.txt.user)

                rm -f /boot/boot.scr.txt.user
            fi

            echo "$ccl" | grep -q 'mxc_hdmi.rgb_quant_range' || ccl="$ccl mxc_hdmi.rgb_quant_range=full"
            echo "$ccl" | grep -q 'mxc_hdmi.enable_3d' || ccl="$ccl mxc_hdmi.enable_3d=enable"
            echo "$ccl" | grep -q 'vpu352=' || ccl="$ccl vpu352=0"
            [ -z "$ccl" ] || sed -i "s%setenv customcmdline.*%$ccl%" /boot/boot.scr.txt

            [ ! -e /boot/boot.scr.txt.xbian ] || rm -f /boot/boot.scr.txt.xbian
            cd /boot; ./mks

            update-rc.d imx6-bluetooth defaults >/dev/null 2>&1
        fi

        INITRD=No run-parts --new-session --report -a $KVER /etc/kernel/postinst.d > /dev/null 2>&1

fi

echo "xbian-package-kernel" >>  /var/run/reboot-required || :

exit 0
