# build/content-tpl/etc/init/frandom.conf
# frandom 
#
# Author: mk01 <development@xbian.org>
#
#

env GENERATOR=hwrng

description     "Loads frandom and changes /dev Xrandom devices accordingly - 01-frandom.rules file is having control over this"

start on started udev

pre-start script

        . /etc/default/xbian-rnd

        if [ $GENERATOR = hwrng ]; then
            modprobe -q bcm2708_rng || modprobe -q bcm2708-rng || :
            [ -e /dev/hwrng ] || exit 1
            /etc/init.d/rng-tools start
        else
            [ -e /dev/frandom -a -e /dev/erandom ] && exit 0
            modprobe -q frandom
        fi
end script

post-start script

        . /etc/default/xbian-rnd

        if [ $GENERATOR = frandom ]; then
            udevadm trigger --action=add --sysname-match=frandom || :
            udevadm settle
            [ -e /dev/frandom -a -e /dev/erandom ] && exit 0
            exit 1
        fi
end script

post-stop script
        . /etc/default/xbian-rnd
        if [ $GENERATOR = hwrng ]; then
            /etc/init.d/rng-tools stop || :
            rmmod bcm2708_rng || :
            exit 0
        fi

        rmmod frandom || :
end script
