# build/imx6/content-tpl/etc/init/failedboot.conf
start on runlevel [2345]

task

exec uboot-env del failedboot || :
