# build/imx6/content-tpl/etc/init/imx6-bluetooth.conf
start on started udevtrigger and stopped udevtrigger

script
    [ $(xbian-arch) = iMX6 ] || exit 0

    if [ "$(cat /sys/devices/soc0/soc.1/2100000.aips-bus/2190000.usdhc/mmc_host/mmc0/mmc0:0001/mmc0:0001:1/device)" = '0x4329' ]; then
        ps ax | grep -v grep | grep -q brcm_patchram_plus || /usr/local/sbin/brcm_patchram_plus  -d --patchram /lib/firmware/brcm/bcm4329.hcd --baudrate 3000000 --use_baudrate_for_download /dev/ttymxc3 --enable_hci &
    elif [ "$(cat /sys/devices/soc0/soc.1/2100000.aips-bus/2190000.usdhc/mmc_host/mmc0/mmc0:0001/mmc0:0001:1/device)" = '0x4330' ]; then
        ps ax | grep -v grep | grep -q brcm_patchram_plus || /usr/local/sbin/brcm_patchram_plus -d --patchram /lib/firmware/brcm/bcm4330.hcd --baudrate 3000000 --use_baudrate_for_download /dev/ttymxc3 --enable_hci --no2bytes --tosleep=50000 &
        ####--enable_lpm 
    fi

    for d in $(ls /sys/class/bluetooth); do
        hciconfig $d up || :
    done

end script
