# build/imx6/content-tpl/etc/init/imx6-reg.conf
start on filesystem

task

script
    echo 0 > /sys/devices/soc0/soc.1/busfreq.13/enable || :
    /usr/local/sbin/devmem2 0x020e0018 w 0xffffffff || :
    exit 0
end script
